# software/deploy/samogrow.service
[Unit]
Description=samogrow AI herb garden controller
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Runs on any always-on Linux machine / VM on the same LAN as the garden's
# Wi-Fi plugs and cameras. Adjust User and paths to your install.
User=samogrow
WorkingDirectory=/opt/samogrow/software
EnvironmentFile=/opt/samogrow/.env
ExecStart=/usr/local/bin/bun run src/main.ts
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
